# openflow-1.3/files/etc/hotplug.d/ofiface/00-netstate
[ ifup = "$ACTION" ] && {
	uci_set_state openflow "$INTERFACE" up 1 
	[ -n "$DEVICE" ] && uci_set_state openflow "$INTERFACE" ifname "$DEVICE"
}
