# openflow-1.3/files/etc/init.d/openflow
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

START=43
STOP=57

boot() {
	ofup -a
}

start() {
	ofup -a
}

restart() {
	stop
	start
}

stop() {
	ofdown -a
}
